# conf.d/packaging/rpm/modbus-binding.service
[Unit]
Description=Redpesk shell service

[Service]
Type=simple
ExecStartPre=mkdir -p /run/user/%U/apis/ws
ExecStart=/usr/bin/afb-binder \
	--name=afb-shell \
	--port=12345 \
	--workdir=/usr/share/spawn-binding \
	--ldpaths=. \
    --roothttp=. \
    --token= \
    --verbose \
    --verbose \
    --verbose \
	--ws-server=unix:/run/user/%U/apis/ws/shell

[Install]
WantedBy=multi-user.target
